# rootfs/etc/services.d/udev/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start udev service
# ==============================================================================
exec udevd