# rootfs/etc/services.d/pulseaudio/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start PulseAudio service
# ==============================================================================
pulseaudio --system -vvv &
PULSE_PID=$!

while [ ! -f /data/internal/cookie_unix ]; do
    bashio::log.debug "Waiting for cookie files"
    sleep 1
done

cp /data/internal/cookie_unix /data/external/cookie_unix

wait $PULSE_PID