# rootfs/etc/services.d/pulseaudio/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start PulseAudio service
# ==============================================================================
exec pulseaudio --system -vvv --disallow-exit --exit-idle-time=-1 --disable-shm