# rootfs/etc/services.d/pulseaudio/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start PulseAudio service
# ==============================================================================

udevd --daemon
udevadm trigger
exec pulseaudio --system -vvv