# rootfs/etc/services.d/pulseaudio/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start PulseAudio service
# ==============================================================================
exec pulseaudio --system -vvv