# rootfs/etc/services.d/pulseaudio/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start PulseAudio service
# ==============================================================================
export PULSE_STATE_PATH="/data/states"
export LD_PRELOAD="/usr/local/lib/libjemalloc.so.2"

exec pulseaudio --system -vvv --disallow-exit --exit-idle-time=-1 --disable-shm
